# debian/truenas-installer.service
[Unit]
Description=TrueNAS Installer

[Service]
ExecStart=/usr/bin/python3 -m truenas_installer --server

[Install]
WantedBy=multi-user.target
